# systemd/raccoon-channel-poster.service
[Unit]
Description=Raccoon Telegram Channel Autoposter
After=network-online.target raccoon-sing-box.service
Wants=network-online.target
Requires=raccoon-sing-box.service

[Service]
Type=simple
User=postingbot
Group=postingbot
WorkingDirectory=/home/postingbot/bots/raccoon-channel-poster
EnvironmentFile=/home/postingbot/bots/raccoon-channel-poster/.env
ExecStart=/home/postingbot/bots/raccoon-channel-poster/.venv/bin/python -m bot.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
